# pds2.conf
start on runlevel [0123456]

# Automatically restart process if crashed
respawn
# Essentially lets upstart know the process will detach itself to the background
expect fork

exec /etc/init/pds1.sh
